# examples/com.hermes.agent-node.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- HERMES Agent Node — ARC-4601 -->
    <!-- Persistent daemon that observes bus, links to gateway, dispatches sub-agents -->
    <!-- Docs: https://github.com/dereyesm/hermes/blob/main/spec/ARC-4601.md -->

    <key>Label</key>
    <string>com.hermes.agent-node</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Replace with your Python path and clan directory -->
        <string>/path/to/python</string>
        <string>-m</string>
        <string>hermes.cli</string>
        <string>daemon</string>
        <string>start</string>
        <string>--foreground</string>
        <string>--dir</string>
        <string>/path/to/clan</string>
    </array>

    <!-- Restart on crash (persistent daemon) -->
    <key>KeepAlive</key>
    <true/>

    <!-- Low priority background process -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/hermes-agent-node-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/hermes-agent-node-stderr.log</string>

    <!-- Environment -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/path/to/hermes/reference/python</string>
    </dict>
</dict>
</plist>
